# docker/8.4/supervisord.conf
[supervisord]
nodaemon=true

[program:horizon]
command=php /var/www/html/artisan horizon
user=sail
autostart=true
autorestart=true
stdout_logfile=/var/www/html/storage/logs/horizon.log
redirect_stderr=true
